# sung.hpc.crontab
# By Sung (23/Oct/2014)
# use the hash sign to prefix a comment
# +---------------- minute (0 - 59)
# |  +------------- hour (0 - 23)
# |  |  +---------- day of month (1 - 31)
# |  |  |  +------- month (1 - 12)
# |  |  |  |  +---- day of week (0 - 7) (Sunday=0 or 7)
# |  |  |  |  |
# *  *  *  *  *  command to be executed
PATH=/home/ssg29/Install/CPAT/usr/local/Cluster-Apps/python/2.7.5/bin:/home/ssg29/Install/bedops_linux_x86_64-v2.4.15/bin:/home/ssg29/Install/bamUtil/bin:/home/ssg29/Install/userApps/bin:/home/ssg29/Install/sratoolkit/bin:/home/ssg29/Install/bwa:/home/ssg29/Install/stringtie:/home/ssg29/Install/miRanda/bin:/home/ssg29/Install/randfold-2.0:/home/ssg29/Install/ViennaRNA/bin:/home/ssg29/Install/bowtie:/home/ssg29/Install/mirdeep2_0_0_7:/home/ssg29/Install/STAR_2.4.2a/bin/Linux_x86_64_static:/home/ssg29/Install/cutadapt-1.8.1/bin:/home/ssg29/Install/trim_galore_v0.4.0:/home/ssg29/Install/fastqc_v0.11.4:/home/ssg29/Install/htslib:/home/ssg29/Install/bcftools:/home/ssg29/Install/vcftools/bin:/home/ssg29/Install/MACS2/bin:/home/ssg29/Install/LiftOver:/home/ssg29/Install/GCcount:/home/ssg29/Install/FREEC:/home/ssg29/Install/ngsutils/bin:/home/ssg29/.local/bin:/home/ssg29/bin:/home/ssg29/Install/samtools:/home/ssg29/Install/methpipe/bin:/usr/local/Cluster-Apps/java/jdk1.6.0_24/bin:/usr/local/Cluster-Apps/parallel/20131222/bin:/usr/local/Cluster-Apps/htop/sbin:/usr/local/Cluster-Apps/htop/bin:/usr/local/Cluster-Apps/R/3.1.1/bin:/usr/local/Cluster-Apps/perl/5.20.0/bin:/usr/local/Cluster-Apps/bedtools/2.20.1/bin:/usr/local/Cluster-Apps/cufflinks/2.2.1/bin:/usr/local/Cluster-Apps/tophat/2.0.12/bin:/usr/local/Cluster-Apps/bismark/0.12.3:/usr/local/Cluster-Apps/bowtie/2.2.3:/usr/local/Cluster-Apps/python/2.7.5/wx/2.9.5.0/bin:/usr/local/Cluster-Apps/python/2.7.5/bin:/usr/local/Cluster-Apps/vtk/5.6.1/bin:/home/ssg29/Install/CPAT/usr/local/Cluster-Apps/python/2.7.5/bin:/home/ssg29/Install/bedops_linux_x86_64-v2.4.15/bin:/home/ssg29/Install/bamUtil/bin:/home/ssg29/Install/userApps/bin:/home/ssg29/Install/sratoolkit/bin:/home/ssg29/Install/bwa:/home/ssg29/Install/stringtie:/home/ssg29/Install/miRanda/bin:/home/ssg29/Install/randfold-2.0:/home/ssg29/Install/ViennaRNA/bin:/home/ssg29/Install/bowtie:/home/ssg29/Install/mirdeep2_0_0_7:/home/ssg29/Install/STAR_2.4.2a/bin/Linux_x86_64_static:/home/ssg29/Install/cutadapt-1.8.1/bin:/home/ssg29/Install/trim_galore_v0.4.0:/home/ssg29/Install/fastqc_v0.11.4:/home/ssg29/Install/htslib:/home/ssg29/Install/bcftools:/home/ssg29/Install/vcftools/bin:/home/ssg29/Install/MACS2/bin:/home/ssg29/Install/LiftOver:/home/ssg29/Install/GCcount:/home/ssg29/Install/FREEC:/home/ssg29/Install/ngsutils/bin:/home/ssg29/.local/bin:/home/ssg29/bin:/home/ssg29/Install/samtools:/home/ssg29/Install/methpipe/bin:/usr/local/Cluster-Apps/intel/fce/12.1.10.319/bin:/usr/local/Cluster-Apps/intel/cce/12.1.10.319/bin:/usr/local/Cluster-Apps/global/bin:/usr/local/Cluster-Apps/intel/impi/4.1.3.045/bin64:/usr/local/Cluster-Apps/vgl/2.3.1/64/bin:/usr/local/Cluster-Apps/turbovnc/2.0.1/bin:/usr/local/Cluster-Apps/java/jdk1.8.0_45/bin:/usr/local/Cluster-Apps/slurm/sbin:/usr/local/Cluster-Apps/slurm/bin:/usr/lib64/qt-3.3/bin:/usr/local/bin:/bin:/usr/bin:/usr/local/sbin:/usr/sbin:/sbin:/opt/ibutils/bin:.:/opt/dell/srvadmin/bin:/home/ssg29/bin

##########
# BACKUP #
##########
0 0 * * * crontab -l > /home/ssg29/config/sung.hpc.crontab
0 0 * * 0 /usr/local/Cluster-Apps/slurm/bin/mybalance >> /home/ssg29/Pipelines/mybalance.txt; /bin/date >> /home/ssg29/Pipelines/mybalance.txt
0 0 * * 0 /usr/bin/du -h --max-depth=2 /scratch/ssg29 >> /home/ssg29/Pipelines/du.txt; /bin/date >> /home/ssg29/Pipelines/du.txt
0 0 * * 0 /usr/bin/du -h --max-depth=2 /scratch/obsgynae/POPS >> /scratch/obsgynae/POPS/du.txt ; /bin/date >> /scratch/obsgynae/POPS/du.txt
0 0 * * 0 /usr/bin/du -h --max-depth=1 /scratch/obsgynae/POPS/data/fastq >> /scratch/obsgynae/POPS/du.fastq.txt ; /bin/date >> /scratch/obsgynae/POPS/du.fastq.txt
0 0 * * 0 /usr/bin/du -h --max-depth=1 /scratch/obsgynae/POPS/results >> /scratch/obsgynae/POPS/du.results.txt ; /bin/date >> /scratch/obsgynae/POPS/du.results.txt
0 0 * * * rsync -a --update /scratch/ssg29/data/fastq /scratch/obsgynae/POPS/data/
#0 1 * * 0 rsync -a --update /scratch/ssg29/data/RoadMap /scratch/obsgynae/POPS/data/
#0 0 * * 0 rsync -a --update /scratch/ssg29/data /scratch/obsgynae/POPS/
0 0 1,15 * * rsync -a --update /scratch/ssg29/results /scratch/obsgynae/POPS/
