# crontab/sung.peta4.crontab
# By Sung (23/Oct/2014)
# use the hash sign to prefix a comment
# +---------------- minute (0 - 59)
# |  +------------- hour (0 - 23)
# |  |  +---------- day of month (1 - 31)
# |  |  |  +------- month (1 - 12)
# |  |  |  |  +---- day of week (0 - 7) (Sunday=0 or 7)
# |  |  |  |  |
# *  *  *  *  *  command to be executed
PATH=:/home/ssg29/Install/gffcompare:/home/ssg29/Install/gffread:/home/ssg29/Install/pigz-2.4:/home/ssg29/Install/bedtools/bin:/home/ssg29/Install/bedops_linux_x86_64-v2.4.15/bin:/home/ssg29/Install/methpipe/bin:/home/ssg29/Install/bamUtil/bin:/home/ssg29/Install/userApps/bin:/home/ssg29/Install/miRanda/bin:/home/ssg29/Install/randfold-2.0:/home/ssg29/Install/ViennaRNA/bin:/home/ssg29/Install/mirdeep2_0_0_7:/home/ssg29/Install/MACS2/bin:/home/ssg29/Install/LiftOver:/home/ssg29/Install/GCcount:/home/ssg29/Install/FREEC:/home/ssg29/Install/ngsutils/bin:/home/ssg29/Install/bcftools:/home/ssg29/Install/vcftools/bin:/home/ssg29/Install/sratoolkit/bin:/home/ssg29/Install/Salmon/bin:/home/ssg29/Install/subread/bin:/home/ssg29/Install/stringtie:/home/ssg29/Install/STAR_2.4.2a/bin/Linux_x86_64_static:/home/ssg29/Install/hisat2:/home/ssg29/Install/bwa:/home/ssg29/Install/bowtie:/home/ssg29/Install/trim_galore_v0.4.0:/home/ssg29/Install/fastqc_v0.11.4:/home/ssg29/Install/htslib:/home/ssg29/Install/samtools:/home/ssg29/perl5/perlbrew/bin:/home/ssg29/local/bin:/home/ssg29/.local/bin:/home/ssg29/bin:/usr/local/Cluster-Apps/parallel/20131222/bin:/usr/local/Cluster-Apps/htop/sbin:/usr/local/Cluster-Apps/htop/bin:/usr/local/Cluster-Apps/texlive/2015/bin/x86_64-linux:/usr/local/software/pandoc/2.0.6/bin:/usr/local/software/ghc/8.2.2/bin:/usr/local/Cluster-Apps/rstudio/1.1.383/bin:/usr/local/software/spack/develop-290617/spack/opt/spack/linux-rhel7-x86_64/gcc-5.4.0/qt-5.4.2-hcaxhn2iih6rdj46wtqhtwh7fqb7pbdd/bin:/usr/local/software/spack/spack-0.11.2/opt/spack/linux-rhel7-x86_64/gcc-4.8.5/gcc-5.4.0-fis24ggupugiobii56fesif2y3qulpdr/bin:/usr/local/software/spack/develop-290617/spack/opt/spack/linux-rhel7-x86_64/gcc-5.4.0/r-3.4.1-jubrpyn77obfqgoxcctim4sa2wkcfar6/bin:/usr/local/Cluster-Apps/cufflinks/2.2.1/bin:/usr/local/software/spack/develop-290617/spack/var/spack/stage/tophat-2.1.1-qghszl7lujculjdhur7nlymijft3ii2v/bin:/usr/local/software/spack/develop-290617/spack/opt/spack/linux-rhel7-x86_64/gcc-5.4.0/bowtie2-2.3.4.1-td2qanwdantsvc4dwlc7mtp5vwnbejbf/bin:/usr/local/software/spack/spack-0.11.2/opt/spack/linux-rhel7-x86_64/gcc-5.4.0/python-3.6.1-xk7ym4l5glcf6ond7yszv2i5gz3wnv2b/bin:/home/ssg29/Install/gffcompare:/home/ssg29/Install/gffread:/home/ssg29/Install/pigz-2.4:/home/ssg29/Install/bedtools/bin:/home/ssg29/Install/bedops_linux_x86_64-v2.4.15/bin:/home/ssg29/Install/methpipe/bin:/home/ssg29/Install/bamUtil/bin:/home/ssg29/Install/userApps/bin:/home/ssg29/Install/miRanda/bin:/home/ssg29/Install/randfold-2.0:/home/ssg29/Install/ViennaRNA/bin:/home/ssg29/Install/mirdeep2_0_0_7:/home/ssg29/Install/MACS2/bin:/home/ssg29/Install/LiftOver:/home/ssg29/Install/GCcount:/home/ssg29/Install/FREEC:/home/ssg29/Install/ngsutils/bin:/home/ssg29/Install/bcftools:/home/ssg29/Install/vcftools/bin:/home/ssg29/Install/sratoolkit/bin:/home/ssg29/Install/Salmon/bin:/home/ssg29/Install/subread/bin:/home/ssg29/Install/stringtie:/home/ssg29/Install/STAR_2.4.2a/bin/Linux_x86_64_static:/home/ssg29/Install/hisat2:/home/ssg29/Install/bwa:/home/ssg29/Install/bowtie:/home/ssg29/Install/trim_galore_v0.4.0:/home/ssg29/Install/fastqc_v0.11.4:/home/ssg29/Install/htslib:/home/ssg29/Install/samtools:/home/ssg29/perl5/perlbrew/bin:/home/ssg29/local/bin:/home/ssg29/.local/bin:/home/ssg29/bin:/usr/local/Cluster-Apps/intel/2017.4/debugger_2017/gdb/intel64_mic/bin:/usr/local/Cluster-Apps/intel/2017.4/compilers_and_libraries_2017.4.196/linux/mpi/intel64/bin:/usr/local/Cluster-Apps/intel/2017.4/compilers_and_libraries_2017.4.196/linux/bin/intel64:/usr/local/software/global/bin:/opt/singularity/bin:/usr/local/Cluster-Apps/singularity/images:/usr/local/Cluster-Apps/vgl/2.5.1/64/bin:/usr/local/Cluster-Apps/turbovnc/2.0.1/bin:/usr/local/software/slurm/current/sbin:/usr/local/software/slurm/current/bin:/usr/lib64/qt-3.3/bin:/home/ssg29/perl5/bin:/usr/local/bin:/usr/bin:/usr/local/sbin:/usr/sbin:/opt/dell/srvadmin/bin:.:/home/ssg29/bin

##########
# BACKUP #
##########
0 0 * * * crontab -l > /home/ssg29/config/sung.peta4.crontab
0 0 * * 0 /usr/local/software/slurm/current/bin/mybalance >> /home/ssg29/Pipelines/mybalance.peta4.txt; /bin/date >> /home/ssg29/Pipelines/mybalance.peta4.txt
0 0 * * 0 /usr/bin/du -h --max-depth=2 /home/ssg29/rds/rds-ssg29-obsgynae/POPS >> /home/ssg29/Pipelines/du.rds.txt; /bin/date >> /home/ssg29/Pipelines/du.rds.txt
0 0 * * 0 /usr/bin/du -h --max-depth=2 /home/ssg29/rcs/rcs-ssg29-obsgynae/POPS >> /home/ssg29/rcs/rcs-ssg29-obsgynae/POPS/du.rcs.txt; /bin/date >> /home/ssg29/rcs/rcs-ssg29-obsgynae/POPS/du.rcs.txt
0 0 * * 0 /usr/bin/du -h --max-depth=1 /home/ssg29/rcs/rcs-ssg29-obsgynae/POPS/data/fastq >> /home/ssg29/rcs/rcs-ssg29-obsgynae/POPS/du.rcs.fastq.txt; /bin/date >> /home/ssg29/rcs/rcs-ssg29-obsgynae/POPS/du.rcs.fastq.txt
0 0 * * 0 /usr/bin/du -h --max-depth=1 /home/ssg29/rcs/rcs-ssg29-obsgynae/POPS/results >> /home/ssg29/rcs/rcs-ssg29-obsgynae/POPS/du.rcs.results.txt; /bin/date >> /home/ssg29/rcs/rcs-ssg29-obsgynae/POPS/du.rcs.results.txt
#0 0 * * * rsync -a --update /home/ssg29/rds/rds-ssg29-obsgynae/POPS/data/fastq /home/ssg29/rcs/rcs-ssg29-obsgynae/POPS/data/
##0 0 * * 0 rsync -a --update /home/ssg29/rds/rds-ssg29-obsgynae/POPS/data /home/ssg29/rcs/rcs-ssg29-obsgynae/POPS/
#0 0 1,15 * * rsync -a --update /home/ssg29/rds/rds-ssg29-obsgynae/POPS/results /home/ssg29/rcs/rcs-ssg29-obsgynae/POPS/
